# .container/crontab
*  *  *  *  *    cd /app && yarn cmd:fetch-locations
